# infra/systemd/engram-obsidian-export.timer
# Timer para exportar Engram → Obsidian cada 5 minutos
# Install: sudo ln -s ...timer /etc/systemd/system/ && sudo systemctl enable --now engram-obsidian-export.timer

[Unit]
Description=Engram → Obsidian export (every 5 min)
Requires=engram-obsidian-export.service

[Timer]
OnCalendar=*:0/5
Persistent=true

[Install]
WantedBy=timers.target
