# infra/systemd/sonora-evolution.service
# Sonora Evolution — daily improvement cycle
# Install: sudo ln -s /home/mystic/sonora-digital-corp/infra/systemd/sonora-evolution.service /etc/systemd/system/
# Enable:  sudo systemctl enable --now sonora-evolution.timer

[Unit]
Description=Sonora Evolution — daily improvement cycle (pattern detection, proposal generation, auto-apply)
Documentation=https://github.com/sonoradigitalcorp-H/sonora-digital-corp
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/env PYTHONPATH=/home/mystic/sonora-digital-corp /usr/bin/python3 -m evolution.auto_evolve --daily
WorkingDirectory=/home/mystic/sonora-digital-corp
User=mystic
Group=mystic

# Security hardening
NoNewPrivileges=yes
ProtectHome=read-only
ProtectSystem=full
PrivateTmp=yes

[Install]
WantedBy=default.target
