# infra/systemd/whatsapp-sync.service
# WhatsApp Sync — Mantiene la sesión de WhatsApp sincronizada en tiempo real
# Corre wacli sync en background para recibir mensajes nuevos
# Install: ln -s /home/ubuntu/sonora-digital-corp/infra/systemd/whatsapp-sync.service /home/ubuntu/.config/systemd/user/
# Enable:  systemctl --user enable --now whatsapp-sync

[Unit]
Description=WhatsApp Sync — real-time message sync daemon
Documentation=https://github.com/sonoradigitalcorp-H/Sonora-Digital-Corp
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/wacli sync --store /home/ubuntu/.config/ai.opencode.desktop/wacli
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/ubuntu/sonora-digital-corp/state/logs/whatsapp-sync.log
StandardError=append:/home/ubuntu/sonora-digital-corp/state/logs/whatsapp-sync.log

[Install]
WantedBy=default.target
