# infra/abe-service.service
[Unit]
Description=ABE Music OS — Internal channel, voice-first label OS
After=network.target docker.service neo4j.service
Wants=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/sonora-digital-corp
Environment=PYTHONPATH=/home/ubuntu/sonora-digital-corp:/home/ubuntu/sonora-digital-corp/apps:/home/ubuntu/sonora-digital-corp/apps/jarvis/src
ExecStart=/home/ubuntu/hermes-agent/venv/bin/python -m uvicorn apps.abe_service.main:app --host 127.0.0.1 --port 5180 --log-level info
Restart=always
RestartSec=5
StandardOutput=append:/home/ubuntu/sonora-digital-corp/logs/abe-service.log
StandardError=append:/home/ubuntu/sonora-digital-corp/logs/abe-service.log

[Install]
WantedBy=multi-user.target
